# memento/sprout.supervisord.conf
[program:memento]
command=/etc/init.d/sprout run
autostart=false
autorestart=true
